# warp/templates/default/warp_upstart.conf.erb
description "Start the twistd process that serves warp webapps"
author "Hoang Xuan Phu <phu@cogini.com>"

start on runlevel [2345]
stop on runlevel [06]

exec twistd \
    -u`id -u vagrant` \
    -g`id -g vagrant` \
    --pidfile=/home/vagrant/twistd.pid \
    --logfile=/home/vagrant/twistd.log \
    -n warp \
    --siteDir=<%= @site_dir %>

respawn
